# rtperfui.container
[Unit]
Description=rtperfui real-time test UI container

[Container]
Image=docker.io/insatomcat/rtperfui:latest
ContainerName=rtperfui
PublishPort=8000:8000
User=root

# Capabilities / ulimits for real-time
PodmanArgs=--privileged
AddCapability=CAP_SYS_NICE
Ulimit=rtprio=99
Ulimit=memlock=-1

# Expose host kernel/sysfs needed for checks and RT tools
# Note: /proc is NOT mounted — in a privileged container the container's own /proc
# already reflects host kernel values (cmdline, cpuinfo, meminfo, sys/*, irq, etc.)
Volume=/sys:/sys
Volume=/run/tuned:/run/tuned:ro
Volume=/etc/tuned:/etc/tuned:ro
Volume=/dev/cpu_dma_latency:/dev/cpu_dma_latency

# Runtime directories (covers pacemaker variants under /run/*)
Volume=/run:/run:ro
# Pacemaker state (CIB + pengine inputs)
Volume=/var/lib/pacemaker:/var/lib/pacemaker:ro

# Corosync config — used to detect cluster mode
Volume=/etc/corosync/corosync.conf:/etc/corosync/corosync.conf:ro

# Ceph config + keyring — needed for rbd image-meta get (VM XML retrieval)
Volume=/etc/ceph:/etc/ceph:ro

Environment=PYTHONUNBUFFERED=1
[Service]
Restart=always

[Install]
WantedBy=default.target